# deploy/com.unarylab.peon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.unarylab.peon</string>
  <!-- caffeinate -i holds off idle/maintenance sleep while peon runs, so a long
       agent run is not cut off when a laptop idles on battery. Closing the lid on
       battery still sleeps; for true always-on, use a non-sleeping host. -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string><string>-i</string><string>--</string>
    <string>/usr/bin/env</string>
    <string>conda</string><string>run</string><string>-n</string><string>peon</string>
    <string>python</string><string>-m</string><string>src</string>
  </array>
  <!-- EDIT: set WorkingDirectory to your local checkout path. -->
  <key>WorkingDirectory</key><string>/Users/YOU/Projects/peon</string>
  <!-- Runtime config (OBSIDIAN_VAULT_PATH, secrets, timeouts, toggles) lives in
       .env at WorkingDirectory; peon loads it first and the CLI subprocess
       inherits it. Set a var here only if it must exist before .env is read
       (e.g. PATH so launchd can find `claude`/`codex` in ~/.local/bin). -->
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <!-- Wait 30s between relaunches so a crash into a still-broken network does not
       hot-loop (default is 10s). -->
  <key>ThrottleInterval</key><integer>30</integer>
</dict>
</plist>
